# luci/applications/luci-youku/root/etc/uci-defaults/luci-youku
#!/bin/sh

uci -q batch <<-EOF >/dev/null
	delete ucitrack.@youku[-1]
	add ucitrack youku
	set ucitrack.@youku[-1].init=youku
	commit ucitrack
EOF

rm -f /tmp/luci-indexcache
exit 0
