# luci/applications/luci-ip-qos/root/etc/init.d/ip-qos
#!/bin/sh /etc/rc.common

START=99

start(){
/usr/bin/ip-qos start
}

stop(){
/usr/bin/ip-qos stop
}
restart(){
/usr/bin/ip-qos stop
/usr/bin/ip-qos start
}



